# httpd/assets/etc/systemd/system/sockets.target.wants/systemd-journal-remote.socket
[Unit]
description=Journal Remote Sink Socket
[Socket]
ListenStream=19532
[Install]
WantedBy=sockets.target 
